# Node_One/filecheck.path
[Unit]
Description= --------Check for files -----------------------------

[Path]
PathExistsGlob=/home/pi/readings/hmacs/hmacreadings*.txt


[Service]
Restart=always
RestartSec=300

[Install]
WantedBy=default.target
